# freenetdev.service
[Unit]
Description=Free netdev

[Service]
ExecStart=/usr/local/bin/auto-free-netdev.sh

[Install]
WantedBy=multi-user.target
